# rust-ananicy.service
[Unit]
Description=Auto Nice Daemon (Rust)
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/rust-ananicy start
Restart=always
# Avoid requiring unit-level nice lowering; daemon adjusts per-process
# CPUAccounting/MemoryAccounting deprecated in recent systemd; removed
CapabilityBoundingSet=CAP_SYS_NICE CAP_SYS_RESOURCE CAP_KILL
AmbientCapabilities=CAP_SYS_NICE CAP_SYS_RESOURCE CAP_KILL
NoNewPrivileges=false
LimitNICE=-20
Delegate=yes

[Install]
WantedBy=multi-user.target